# matrix.service
[Unit]
Description=Matrix Jump Protocol Listener
After=network.target

[Service]
Type=simple
WorkingDirectory=/root/Matrix
ExecStart=/usr/local/bin/matrix listen
EnvironmentFile=-/root/Matrix/.env
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
